# .config/systemd/user/rclone-sync@songbook.timer
[Unit]
Description=Timer para rclone-bisync (songbook)

[Timer]
OnBootSec=5m
OnUnitActiveSec=1m
Persistent=true

[Install]
WantedBy=timers.target
