# .config/systemd/user/rclone-sync@.service
[Unit]
Description=Rclone bisync com perfil %i

[Service]
Type=oneshot
ExecStart=/home/gioz/.local/bin/rclone-sync.sh %i
